# env2host.service
[Unit]
Description=Kubernetes env to hosts vars

[Service]
Type=oneshot
ExecStart=/usr/local/bin/env2host.php

[Install]
WantedBy=multi-user.target